# lprint.service.in
[Unit]
Description=LPrint Service
After=network.target nss-lookup.target
Requires=avahi-daemon.socket

[Install]
WantedBy=multi-user.target

[Service]
ExecStart=@bindir@/lprint server -o log-file=- -o log-level=info
Type=simple
Restart=on-failure

