# lprint.service.in
[Unit]
Description=LPrint Service
After=network.target nss-lookup.target
Requires=avahi-daemon.socket

[Install]
WantedBy=multi-user.target

[Service]
ExecStart=@bindir@/lprint server -o admin-group=@ADMINGROUP@ -o auth-service=@AUTHSERV@ @LISTENHOST@-o log-file=@LOGFILE@ -o log-level=info -o server-options=@SOPTIONS@
Type=simple
Restart=on-failure

